# stopdb/redbus.timer
[Unit]
Description=Daily redbus update

[Timer]
OnCalendar=daily

[Install]
WantedBy=multi-user.target
